# deployment/configs/backend.service
# deployment/configs/backend.service
[Unit]
Description=CouldYou Backend Service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/couldyou/backend
Environment="PATH=/opt/couldyou/backend/venv/bin"
ExecStart=/opt/couldyou/backend/venv/bin/gunicorn --workers 4 --bind 127.0.0.1:5000 run:app
Restart=always

[Install]
WantedBy=multi-user.target
