# deployment/configs/telegram.service
# deployment/configs/telegram.service
[Unit]
Description=CouldYou Telegram Bot Service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/couldyou_chatbot/telegram_bot
Environment="PATH=/opt/couldyou_chatbot/telegram_bot/venv/bin"
EnvironmentFile=/opt/couldyou_chatbot/.env
ExecStart=/opt/couldyou_chatbot/telegram_bot/venv/bin/python bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
